# socket-activation/templates/EXAMPLE.socket
[Unit]
Description=Socket for EXAMPLE


[Socket]

ListenStream=GLOBAL_IP:GLOBAL_PORT
Transparent=True

[Install]
WantedBy=sockets.target